# config/autostart/xscreensaver-properties.desktop
[Desktop Entry]
Exec=xscreensaver -no-splash
Icon=xscreensaver
Terminal=false
Name=Screensaver
Type=Application
